# script/supervisord.ini
[supervisord]

[program:open-admin-server]
directory=/usr/local/services/sql-operator/
command=/usr/local/services/sql-operator/sql-operator
autostart=true
autorestart=true
startsecs=10
startretries=10
stopasgroup=true
killasgroup=true
redirect_stderr=true
stdout_logfile=/app/logs/console.log
stdout_logfile_maxbytes=200MB
stdout_logfile_backups=20
